# GoRepost.service
[Unit]
Description=GoRepost daemon
After=networking.service
 
[Service]
Restart=on-failure
ExecStart=/path/to/binnary/GoRepost-latest
